# cv6/ws-js.conf
[program:ws-js]
directory=/home/xsefcikj/chatserver-js
command=/home/xsefcikj/.nvm/versions/node/v22.14.0/bin/node index.js
autostart=true
autorestart=true
stderr_logfile=/home/xsefcikj/chatserver-js/server.err.log
stdout_logfile=/home/xsefcikj/chatserver-js/server.out.log
